# yshvydak-dashboard.service
[Unit]
Description=YShvydak Test Dashboard - Production Deployment
Documentation=https://github.com/yshvydak/yshvydak-test-dashboard
After=network.target

[Service]
Type=forking
User=pi
WorkingDirectory=/home/pi/yshvydak-test-dashboard
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/pm2 start /home/pi/yshvydak-test-dashboard/ecosystem.config.js
ExecReload=/usr/local/bin/pm2 restart /home/pi/yshvydak-test-dashboard/ecosystem.config.js
ExecStop=/usr/local/bin/pm2 stop /home/pi/yshvydak-test-dashboard/ecosystem.config.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yshvydak-dashboard

[Install]
WantedBy=multi-user.target
